# home/.config/profile.d/40-editor
# Use neovide as the default editor in GUI sessions, nvim otherwise.
if [ -n "$DISPLAY" ]; then
  EDITOR=neovide
else
  EDITOR=nvim
fi

VISUAL=$EDITOR
