# home/.config/profile.d/60-XCompose
export XCOMPOSEFILE=~/.config/X11/Compose
export GTK_IM_MODULE="xim"
